# launchd/com.lyeyixian.rental-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Fetch agent: opens a visible Chromium window daily at 19:00 (and on login),
  calls Mastercard's conversion-rates API for the 1st of the current month,
  and caches the result in local/state.json. The script's self-dedup (see
  src/state.ts) makes every firing after the first successful fetch of the
  month a millisecond no-op; the script also exits silently on the 1st of
  the month (PRD #7's "from the 2nd" floor). See ADR-0007 for why the trigger
  is daily-calendar rather than login-only.

  Tokens to replace before installing into ~/Library/LaunchAgents/:

    __REPO_PATH__   absolute path to this repository checkout
                    (e.g. /Users/yourname/Documents/repo/rental-mastercard-calculator)
    __PNPM__        absolute path to the pnpm binary (find it with `command -v pnpm`;
                    e.g. /opt/homebrew/bin/pnpm on Apple Silicon, /usr/local/bin/pnpm
                    on Intel). launchd needs an absolute path here — unlike a shell it
                    does not search PATH to resolve the command name.
    __PATH__        PATH the agent should inherit. Must include the directory holding
                    `node`, which pnpm shells out to. A typical Apple-Silicon value:
                    /opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
                    On Intel macOS, use:
                    /usr/local/bin:/usr/bin:/bin
    __LOG_DIR__     directory for the agent's stdout/stderr log, e.g.
                    /Users/yourname/Library/Logs. Must be OUTSIDE the
                    TCC-protected folders (~/Documents, ~/Desktop, ~/Downloads):
                    launchd opens the log file before spawning the job, macOS
                    privacy protection denies that open for a non-Apple program,
                    and the job dies with a silent exit 78 (EX_CONFIG) and an
                    empty log. See ADR-0009.

  See README "Running autonomously with launchd" for install / verify / uninstall.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lyeyixian.rental-fetch</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PNPM__</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_PATH__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/rental-fetch.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/rental-fetch.log</string>
</dict>
</plist>
